# msyncd/bin/msyncd.service
[Unit]
Description=Sync FW daemon
Requires=dbus.socket

[Service]
ExecStart=/usr/bin/msyncd
